# etc/systemd/system/torb.python.socket
[Unit]
Description=torb.python socket

[Socket]
ListenStream=/home/isucon/run/torb.python/socket

[Install]
WantedBy=sockets.target
